# bin/devcontestd
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

CONFIG_FILE="/usr/share/devcontest/config.ini"

depend() {
	need net
	after sshd
}

start() {
	paster serve --daemon ${CONFIG_FILE}
}

stop() {
	paster serve --stop-daemon ${CONFIG_FILE}
}

restart() {
	svc_stop; svc_start;
}

